# asset/systemd/auto-clear-log.service
[Unit]
Description=Fix network
After=network.target

[Service]
Type=simple
ExecStart=/bin/auto-clear-log

[Install]
WantedBy=network.target